# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for nathangould.com"

start on net-device-up
stop on shutdown

respawn

setuid nathan
chdir /home/nathan/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    superlists.wsgi:application
